# configuration/systemd/noise.service
# Service file for noise

[Unit]
Description=Noise detection
After=network.target syslog.target motion.service

[Service]
Type=simple
ExecStart=/opt/fruitnanny/bin/noise-daemon.py
LimitNOFILE=8000
WorkingDirectory=/opt/fruitnanny
KillMode=mixed
Restart=on-failure
User=pi
Group=pi
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=noise
SuccessExitStatus=SIGKILL

[Install]
WantedBy=multi-user.target
